# sphttp.conf
[group:hdpt] 
programs=server8080,server8081,server8082,server8083,server8084,server8085,server8086,server8087,server8088,server8089,server8090,autoclick

[program:server8080]
command=python /home/DOme/dome_server/server.py --port=8080  
directory=/home
user=root
autostart=true      
startsecs=10        
autorestart=true     
startretries=3       
redirect_stderr=true 
stdout_logfile_maxbytes=20MB  
stdout_logfile_backups = 20   
stdout_logfile=home
stopasgroup=true 
killasgroup=true 
loglevel=info

[program:server8081]
command=python /home/DOme/dome_server/server.py --port=8081  
directory=/home
user=root
autostart=true       
startsecs=10        
autorestart=true    
startretries=3       
redirect_stderr=true 
stdout_logfile_maxbytes=20MB  
stdout_logfile_backups = 20   
stdout_logfile=home
stopasgroup=true 
killasgroup=true 
loglevel=info


[program:server8082]
command=python /home/DOme/wxServer/server.py --port=8082  
directory=/home
user=root
autostart=true      
startsecs=10        
autorestart=true    
startretries=3      
redirect_stderr=true 
stdout_logfile_maxbytes=20MB  
stdout_logfile_backups = 20   
stdout_logfile=home
stopasgroup=true 
killasgroup=true 
loglevel=info

[program:server8083]
command=python /home/DOme/wxServer/server.py --port=8083 
directory=/home
user=root
autostart=true       
startsecs=10        
autorestart=true     
startretries=3      
redirect_stderr=true 
stdout_logfile_maxbytes=20MB  
stdout_logfile_backups = 20   
stdout_logfile=home
stopasgroup=true 
killasgroup=true 
loglevel=info

[program:server8083]
command=python /home/DOme/wxServer/server.py --port=8083  
directory=/home
user=root
autostart=true       
startsecs=10        
autorestart=true     
startretries=3       
redirect_stderr=true 
stdout_logfile_maxbytes=20MB 
stdout_logfile_backups = 20   
stdout_logfile=home; ÈÕÖ¾Â·¾¶
stopasgroup=true
killasgroup=true 
loglevel=info

[program:server8083]
command=python /home/DOme/wxServer/server.py --port=8083  
directory=/home
user=root
autostart=true       
startsecs=10        
autorestart=true     
startretries=3       
redirect_stderr=true 
stdout_logfile_maxbytes=20MB  
stdout_logfile_backups = 20   
stdout_logfile=home
stopasgroup=true 
killasgroup=true 
loglevel=info

[program:server8084]
command=python /home/DOme/wxServer/server.py --port=8084  
directory=/home
user=root
autostart=true       
startsecs=10        
autorestart=true    
startretries=3       
redirect_stderr=true 
stdout_logfile_maxbytes=20MB  
stdout_logfile_backups = 20 
stdout_logfile=home
stopasgroup=true 
killasgroup=true 
loglevel=info

program:server8085]
command=python /home/DOme/wxServer/server.py --port=8085
directory=/home
user=root
autostart=true       
startsecs=10        
autorestart=true    
startretries=3       
redirect_stderr=true 
stdout_logfile_maxbytes=20MB  
stdout_logfile_backups = 20 
stdout_logfile=home
stopasgroup=true 
killasgroup=true 
loglevel=info

program:server8086]
command=python /home/DOme/wxServer/server.py --port=8086
directory=/home
user=root
autostart=true       
startsecs=10        
autorestart=true    
startretries=3       
redirect_stderr=true 
stdout_logfile_maxbytes=20MB  
stdout_logfile_backups = 20 
stdout_logfile=home
stopasgroup=true 
killasgroup=true 
loglevel=info

program:server8087]
command=python /home/DOme/wxServer/server.py --port=8087
directory=/home
user=root
autostart=true       
startsecs=10        
autorestart=true    
startretries=3       
redirect_stderr=true 
stdout_logfile_maxbytes=20MB  
stdout_logfile_backups = 20 
stdout_logfile=home
stopasgroup=true 
killasgroup=true 
loglevel=info

program:server8088]
command=python /home/DOme/wxServer/server.py --port=8088
directory=/home
user=root
autostart=true       
startsecs=10        
autorestart=true    
startretries=3       
redirect_stderr=true 
stdout_logfile_maxbytes=20MB  
stdout_logfile_backups = 20 
stdout_logfile=home
stopasgroup=true 
killasgroup=true 
loglevel=info




program:server8089]
command=python /home/DOme/wxServer/server.py --port=8089
directory=/home
user=root
autostart=true       
startsecs=10        
autorestart=true    
startretries=3       
redirect_stderr=true 
stdout_logfile_maxbytes=20MB  
stdout_logfile_backups = 20 
stdout_logfile=home
stopasgroup=true 
killasgroup=true 
loglevel=info

program:server8090]
command=python /home/DOme/wxServer/server.py --port=8090
directory=/home
user=root
autostart=true       
startsecs=10        
autorestart=true    
startretries=3       
redirect_stderr=true 
stdout_logfile_maxbytes=20MB  
stdout_logfile_backups = 20 
stdout_logfile=home
stopasgroup=true 
killasgroup=true 
loglevel=info

[program:autoclick]
command=python /home/DOme/autoDome/auto_click.py 
directory=/home
user=root
autostart=true      
startsecs=10        
autorestart=true     
startretries=3       
redirect_stderr=true 
stdout_logfile_maxbytes=20MB  
stdout_logfile_backups = 20   
stdout_logfile=home
stopasgroup=true 
killasgroup=true
loglevel=info
